# 10-walm-montana.rules
ACTION=="add", SUBSYSTEM=="block", ENV{ID_FS_TYPE}=="ntfs", PROGRAM="/bin/sleep 1.5", RUN+="/bin/sh -c 'ntfsfix -d /dev/%k &'"
