# 10-walm-montana.rules
ACTION=="add"  RUN+="/bin/ntfsfix -d /dev/%k"
